# meshfox/sysdrv/out/rootfs_uclibc_rv1106/etc/systemd/system/serial-getty@.service.d/override.conf
[Unit]
After=femto-boot-complete.service multi-user.target
Conflicts=getty.target
Before=getty.target

[Service]
ExecStartPre=/bin/sleep 0.5
