# meshfox/sysdrv/out/rootfs_uclibc_rv1106/etc/systemd/system/meshfox-boot-complete.service
[Unit]
Description=MeshFox Boot Complete
DefaultDependencies=no
After=rc-local.service

[Service]
ExecStart=/usr/local/bin/meshfox-boot-complete.sh
Type=oneshot
Restart=no
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=default.target
